# docker/mailserver/chpasswd-supervisor.conf
[program:chpasswd]
command=bash -c 'echo "root:$SSH_ROOT_PASSWORD" | chpasswd'
